# auto_update_cache.crontab
* * * * * /home/artem/WEB/WEB-12/Web/web_park_vk/.venv/bin/python3 /home/artem/WEB/WEB-12/Web/web_park_vk/askme/manage.py update_cache
